# deploy/greed-mobile.service
[Unit]
Description=GREED Mobile API (backend-mobile)
After=network.target

[Service]
Type=simple
User=opc
Group=opc
WorkingDirectory=/opt/greed/backend-mobile
EnvironmentFile=/opt/greed/backend-mobile/.env
ExecStart=/opt/greed/backend-mobile/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=greed-mobile

[Install]
WantedBy=multi-user.target
